# ops/topgoldpicks-backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Daily MongoDB backup at 02:00 UTC, 14-day retention.
# Install with: cp ops/topgoldpicks-backup.cron /etc/cron.d/topgoldpicks-backup
0 2 * * * root /opt/topgoldpicks/scripts/backup.sh >> /var/log/topgoldpicks-backup.log 2>&1
